# env/searchstation.conf
[program:collector]
command=/home/ubuntu/searchstation/run.sh
environment=HOME='/home/ubuntu',USER='ubuntu'
directory=/home/ubuntu/searchstation
autostart=true
autorestart=true
startretries=3
